# Deploy/conf/backup_db.service
[Unit]
Description=DB Backup Service

[Service]
ExecStart=/bin/bash /opt/LegacyPlayersV4/Deploy/backup_db.sh
User=root
